# build/com.waelrabadi.onemark.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.waelrabadi.onemark</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Onemark.app/Contents/MacOS/Onemark</string>
        <string>--start=api</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Applications/Onemark.app/Contents</string>
    <key>KeepAlive</key>
    <false/>
    <key>LaunchOnlyOnce</key>        
    <false/>
    <key>Sockets</key>
    <dict>
        <key>Listeners</key>
        <dict>
            <key>SockServiceName</key>
            <string>32801</string>
            <key>SockType</key>
            <string>stream</string>
            <key>SockFamily</key>
            <string>IPv4</string>
        </dict>
    </dict>
</dict>
</plist>